# Temperature/temperature-login.timer
[Unit]
Description=1-second delay before running temperature adjustment

[Timer]
OnStartupSec=1s
OnUnitActiveSec=0

[Install]
WantedBy=timers.target