# Temperature/temperature.timer
[Unit]
Description=Screen Temperature Adjuster Timer

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
